# s6/dnsmasq/run
#!/bin/sh

exec /usr/sbin/dnsmasq -d -C /etc/dnsmasq.conf
